# examples/com.investment.daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Example launchd agent for macOS that runs fetch_report.py every day
  at 09:00 (local time). Before installing, replace every occurrence of
  REPO_PATH with the absolute path to your cloned repo
  (e.g. /Users/yourname/daily-investment-report) and copy this file to
  ~/Library/LaunchAgents/, then `launchctl load <path>`.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investment.daily-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>REPO_PATH/.venv/bin/python3</string>
        <string>REPO_PATH/fetch_report.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REPO_PATH</string>

    <!-- 每天 09:00（不含 Weekday key 即每天触发）。只想工作日跑时，按 Weekday 1-5 拆成 5 条。 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>REPO_PATH/reports/cron.log</string>

    <key>StandardErrorPath</key>
    <string>REPO_PATH/reports/cron_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
